# systemd/backup-scheduler.service
# Backup scheduler daemon. Runs "Backup backup schedule" in a loop,
# executing backups whose cron schedule matches the current time.
#
# Install:
#   sudo cp backup-scheduler.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable backup-scheduler
#   sudo systemctl start backup-scheduler
#
# Prerequisites:
#   - PostgreSQL with schema iks (run schema/init_iks_schema.sql or start the UI once).
#   - Backup binary at ExecStart path (e.g. /usr/local/bin/Backup).
#   - Configuration via environment only: set POSTGRES_* (or put backup.env in WorkingDirectory).
#   - Binary reads backup.env / .env from WorkingDirectory; or set Environment= below.

[Unit]
Description=Backup scheduler (cron-style backup runner)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/Backup backup schedule
WorkingDirectory=/etc/backup
# Set PostgreSQL connection (or use backup.env in WorkingDirectory)
# Environment=POSTGRES_HOST=localhost POSTGRES_PORT=5432 POSTGRES_DB=postgres POSTGRES_USER=backup POSTGRES_PASSWORD=secret

Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=backup-scheduler

# Run as dedicated user (create with: useradd -r -s /bin/false backup)
# User=backup
# Group=backup

[Install]
WantedBy=multi-user.target
